# calendar/krasovska-calendar.service
[Unit]
Description=Krasovská Calendar ICS Server
After=network.target docker.service
Wants=docker.service

[Service]
Type=notify
User=leos
WorkingDirectory=/home/leos/calendar
ExecStart=/usr/bin/python3 /home/leos/calendar/ics_server.py
Restart=always
RestartSec=5
WatchdogSec=120
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/leos/calendar/.env

[Install]
WantedBy=multi-user.target
